# packages/root_backup/dotfiles/restic-backup.service
[Service]
CPUSchedulingPolicy=idle
Environment=HOME=/root
ExecStart=/usr/bin/restic-backup backup-home

IOSchedulingClass=idle

[Unit]
Description=Restic backup
Wants=network-online.target
Requisite=ac.target
