# packages/rslsync/dotfiles/rslsync.service
[Unit]
Description=Resilio Sync Service
After=network.target
Wants=network-online.target

[Service]
ExecStart=rslsync --nodaemon --config {{@@ env['HOME'] @@}}/.config/rslsync/config
Type=simple
Restart=on-failure
# SuccessExitStatus=3 4
# RestartForceExitStatus=3 4

[Install]
WantedBy=default.target
