# packages/i3/dotfiles/xcape.service
[Service]
ExecStart=xcape -e 'Control_L=Escape'
Type=simple
Restart=on-failure
[Unit]
Description=Xcape
[Install]
WantedBy=graphical-session.target
